# systemd/ecnu-openconnect-keepalive.timer
[Unit]
Description=Periodic keepalive for OpenConnect split-tunnel access

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
RandomizedDelaySec=20s
Unit=ecnu-openconnect-keepalive.service

[Install]
WantedBy=timers.target
