# systemd/och-openconnect-keepalive.service
[Unit]
Description=Keep OpenConnect split-tunnel access alive
After=network-online.target och-openconnect.service
Wants=network-online.target

[Service]
Type=oneshot
Environment=VPN_SERVICE=och-openconnect.service
Environment=CONFIG_FILE=/etc/och/och-vpn.env
Environment=FAIL_THRESHOLD=2
Environment=STATE_DIR=/run/och-openconnect-keepalive
Environment=PING_COUNT=1
Environment=PING_TIMEOUT=1
ExecStart=/usr/local/bin/och-openconnect-keepalive
