# systemd/och-openconnect.service
[Unit]
Description=AnyConnect/OpenConnect split-tunnel access
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
Environment=CONFIG_FILE=/etc/och/och-vpn.env
Environment=PID_FILE=/run/och-openconnect.pid
Environment=LOG_FILE=/var/log/och-openconnect.log
ExecStart=/usr/local/bin/och-vpn connect
ExecStop=/usr/local/bin/och-vpn disconnect
TimeoutStartSec=90
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
